# tools/mini/run.jam.inbox.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>run.jam.inbox</string>
  <!-- The contributor inbox: a folder dropped into mark's own home (a family member over
       Tailscale, rsync'd via his own restricted account — see AGENTS.md §Family/affiliate)
       becomes a station. jam-inbox.sh runs its own 20s poll loop, so KeepAlive (long-running),
       not StartInterval. Needs docker on PATH — the script exports it. INBOX points at
       /Users/mark/jam-inbox deliberately: mark's account owns it outright, so no cross-account
       permission grant is needed for him to write there or for this job (running as jason) to
       read it back out — only default home-dir read perms, nothing special. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/jason/business/jam-station/tools/jam-inbox.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>INBOX</key><string>/Users/mark/jam-inbox</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ThrottleInterval</key><integer>20</integer>
  <key>StandardOutPath</key><string>/Users/jason/Library/Logs/jam-inbox.log</string>
  <key>StandardErrorPath</key><string>/Users/jason/Library/Logs/jam-inbox.log</string>
</dict>
</plist>
