# tools/mini/run.jam.cdwatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>run.jam.cdwatch</string>
  <!-- Program is the jam-cdd helper (compiled from tools/mini/jam-cdd.c). It's the job's
       responsible process, so granting THIS binary Full Disk Access lets the whole watcher —
       cd-tick, rip-cd, lame — read the CD drive on a headless Mac. jam-cdd runs its own 12s
       loop and runs the rip in-job, so KeepAlive (long-running), not StartInterval. -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/jason/bin/jam-cdd</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ThrottleInterval</key><integer>10</integer>
  <key>StandardOutPath</key><string>/dev/null</string>
  <key>StandardErrorPath</key><string>/Users/jason/Library/Logs/jam-cdwatch.log</string>
</dict>
</plist>
