# tools/mini/run.jam.nomusic.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>run.jam.nomusic</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/jason/business/jam-station/tools/mini/no-music.sh</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/jason/Library/Logs/jam-nomusic.log</string>
  <key>StandardErrorPath</key><string>/Users/jason/Library/Logs/jam-nomusic.log</string>
</dict>
</plist>
